# inference.service
[Unit]
Description=TradingBotV4 - Inference Service (Pi)
Documentation=https://github.com/sutefan7/trading-bot-v4
After=network-online.target
Wants=network-online.target
StartLimitInterval=300
StartLimitBurst=3

[Service]
Type=simple
User=stephang
Group=stephang
WorkingDirectory=/home/stephang/trading-bot-v4

# Health check before start
ExecStartPre=/bin/bash -c 'python3 -c "import requests; requests.get(\"https://api.kraken.com/0/public/Time\", timeout=5)"'
ExecStartPre=/bin/bash -c 'test -f /home/stephang/trading-bot-v4/.env'
ExecStartPre=/bin/bash -c 'test -f /home/stephang/trading-bot-v4/latest.txt'
ExecStartPre=/bin/bash -c 'cd /home/stephang/trading-bot-v4 && source venv-light/bin/activate && python smoke_test.py'

# Main service - inference only
ExecStart=/home/stephang/trading-bot-v4/start_inference_bot.sh

# Restart policy with limits
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=3

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/stephang/trading-bot-v4"
Environment="TRADING_BOT_ENV=production"
Environment="TRADING_BOT_MODE=inference"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradingbot-inference

# Resource limits (lighter for inference only)
MemoryMax=1G
CPUQuota=60%
TasksMax=500

# Watchdog
WatchdogSec=60

[Install]
WantedBy=multi-user.target
