# systemd/trading-bot-pi.service
[Unit]
Description=Trading Bot Pi - ML Runner (Pi)
Documentation=https://github.com/sutefan7/trading-bot-pi
After=network-online.target
Wants=network-online.target
StartLimitInterval=300
StartLimitBurst=3

[Service]
Type=simple
User=stephang
Group=stephang
WorkingDirectory=/srv/trading-bot-pi/app

# Health check before start (quote URL correct)
ExecStartPre=/usr/bin/python3 -c 'import requests; requests.get("https://api.kraken.com/0/public/Time", timeout=5); print("Kraken OK")'
ExecStartPre=/bin/bash -c 'test -f /srv/trading-bot-pi/app/storage/artifacts/latest.txt || true'

# Main service (ML-enabled)
ExecStart=/bin/bash -lc 'source /srv/trading-bot-pi/app/.venv/bin/activate && PYTHONPATH=/srv/trading-bot-pi/app python3 src/apps/runner/main_v2_with_ml.py --non-interactive'

# Restart policy with limits
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=3

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/srv/trading-bot-pi/app"
Environment="TRADING_BOT_ENV=production"
Environment="PATH=/srv/trading-bot-pi/app/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot-pi

# Resource limits
MemoryMax=2G
CPUQuota=80%
TasksMax=1000

# Watchdog (disabled to avoid unintended restarts during long operations)
WatchdogSec=0

[Install]
WantedBy=multi-user.target
